# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=bash -c 'sleep 5 && zkServer.sh start-foreground'
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
